# agent/systemd/litellm.service
[Unit]
Description=LiteLLM proxy
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=%h/.config/litellm
Environment=PATH=%h/.local/bin:/run/current-system/sw/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env litellm --config %h/.config/litellm/litellm_config.yaml --host 127.0.0.1 --port 4000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
